# deploy/d-brain-weekly.timer
[Unit]
Description=Run d-brain weekly digest Friday morning

[Timer]
OnCalendar=Fri 06:00
Persistent=true

[Install]
WantedBy=timers.target
