# deploy/d-brain-process.service
[Unit]
Description=d-brain Daily Processing

[Service]
Type=oneshot
User=shima
WorkingDirectory=/home/shima/projects/dbrain
ExecStart=/home/shima/projects/dbrain/scripts/process.sh
Environment=PYTHONUNBUFFERED=1
